# workdirfs.service
[Unit]
Description=Daily working directory ~/Work

[Service]
Environment=WORKDIRFS_OPTIONS="-a %h/archive -m %h/Work -t 3 -M False -y False"
ExecStart=/usr/bin/python3 -u /usr/local/bin/workdirfs.py ${WORKDIRFS_OPTIONS}
ExecStop=/usr/bin/fusermount -u ${DST}
