# workdirfs.service
[Unit]
Description=Daily working directory ~/Work

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/workdirfs.py
