# workdirfs.service
[Unit]
Description=Daily working directory ~/Work

[Service]
Environment=SRC=%h/archive DST=%h/Work
ExecStart=/usr/bin/python3 /usr/local/bin/workdirfs.py ${SRC} ${DST}
ExecStop=/usr/bin/fusermount -u ${DST}
