# workdirfs.service
[Unit]
Description=Daily working directory ~/Work

[Service]
ExecStart=/usr/bin/python3 %h/.local/myshellconfig/bin/workdir.py
